# deploy/systemd/qhat.service
[Unit]
Description=qhat signaling server
After=network.target

[Service]
Type=simple
User=appuser
Group=appuser
WorkingDirectory=/opt/qhat
ExecStart=/opt/qhat/server/bin/qhat-server
EnvironmentFile=-/opt/qhat/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
